# ephesus/cron/crontab
5 3 * * * bash /usr/local/cron/greek-room-keycloak-db-backup.sh -c /usr/local/cron/greek-room-keycloak-db-backup.config
